# installers/darwin/com.claude-relay.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-relay</string>

  <!-- Port and session travel as argv, not env vars. PORT is a generic
       env namespace honored by nearly every Node dev server, so keeping
       it out of the environment makes an entire class of port
       collisions impossible. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{REPO}}/scripts/serve.sh</string>
    <string>--port</string>
    <string>{{PORT}}</string>
    <string>--session</string>
    <string>{{SESSION}}</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{REPO}}</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <!-- Without a UTF-8 locale, tmux falls back to a "safe" mode that
         rewrites control characters (including tabs) to '_' in format
         output, which breaks our list-windows parser. Load-bearing. -->
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>TERM</key>
    <string>xterm-256color</string>
    <key>DEBUG</key>
    <string>{{DEBUG}}</string>
  </dict>

  <!-- Capture stdout/stderr to the repo so first-boot crashes are
       diagnosable. On a healthy daemon these stay effectively empty.
       Set DEBUG=1 in the install environment (DEBUG=1 ./relay install) to
       additionally enable node's rotating app log at
       {{REPO}}/logs/debug.log. -->
  <key>StandardOutPath</key>
  <string>{{REPO}}/logs/launchd.out</string>

  <key>StandardErrorPath</key>
  <string>{{REPO}}/logs/launchd.err</string>
</dict>
</plist>
