# installers/darwin/com.claude-relay.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-relay</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{REPO}}/start.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{REPO}}</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <!-- Without a UTF-8 locale, tmux falls back to a "safe" mode that
         rewrites control characters (including tabs) to '_' in format
         output, which breaks our list-windows parser. Load-bearing. -->
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>TERM</key>
    <string>xterm-256color</string>
    <key>PORT</key>
    <string>{{PORT}}</string>
    <key>TMUX_SESSION</key>
    <string>{{SESSION}}</string>
    <key>DEBUG</key>
    <string>{{DEBUG}}</string>
  </dict>

  <!-- No log files by default. Set DEBUG=1 via --debug at install
       time to enable node-level rotating log file at
       ~/.cache/claude-relay/debug.log (see spec §5.4). -->
  <key>StandardOutPath</key>
  <string>/dev/null</string>

  <key>StandardErrorPath</key>
  <string>/dev/null</string>
</dict>
</plist>
