# installers/linux/claude-relay.service.template
[Unit]
Description=Claude Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Port and session travel as argv, not env vars. PORT is a generic
# env namespace honored by nearly every Node dev server, so keeping
# it out of the environment makes an entire class of port collisions
# impossible.
ExecStart=/bin/bash {{REPO}}/scripts/serve.sh --port {{PORT}} --session {{SESSION}}
WorkingDirectory={{REPO}}
Restart=always
RestartSec=10
StandardOutput=null
StandardError=null
Environment="PATH={{PATH}}"
Environment="HOME={{HOME}}"
Environment="LANG=en_US.UTF-8"
Environment="TERM=xterm-256color"
Environment="DEBUG={{DEBUG}}"

[Install]
WantedBy=default.target
